# systemd/cryptoserver.service
[Unit]
Description=PasswordServer FastAPI (Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
# Adjust these paths to your deployment location
WorkingDirectory=/home/ubuntu/cryptoserver
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/ubuntu/cryptoserver/run.sh
Restart=on-failure
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
