# systemd/passwordserver.service
[Unit]
Description=PasswordServer FastAPI (Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
# Adjust these paths to your deployment location
WorkingDirectory=/opt/passwordserver
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/opt/passwordserver/venv/bin:/usr/bin"
ExecStart=/usr/bin/env uvicorn app:app --host 0.0.0.0 --port 8000 --workers 2
Restart=on-failure
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
